# systemd/netonix-manager.service
[Unit]
Description=Netonix Manager
Documentation=https://github.com/votre-org/netonix-manager
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/netonix-manager
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netonix-manager
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
